# edge-ai.service
[Unit]
Description=Edge AI Inference Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/edge-ai
Environment="QUEUE_SIZE=10"
Environment="PORT=5000"
Environment="HOST=0.0.0.0"
Environment="LOG_LEVEL=INFO"
ExecStart=/home/pi/edge-ai/venv/bin/python3 /home/pi/edge-ai/src/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

